# core/roboconf-target-docker/src/main/resources/rc.local
#!/bin/sh
# Generated by Roboconf

# System update
apt-get update

# Start Roboconf
cd /usr/local/roboconf-agent
./start.sh

# Always return 0
exit 0
